# curiefense/images/confserver/bootstrap/initial-bucket-export.conf
# supervisord configuration file
[program:initial-bucket-export]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0
startsecs=0
command=/bootstrap/initial-bucket-export.sh
autorestart=false
